# backend/crontab
# Graphile Worker crontab
# Format: <cron expression> <task_name> [?options=json]
# Docs: https://worker.graphile.org/docs/cron
#
# expire_needs: mark campaign-needs past their target date as expired.
# Runs once per hour at minute 0. Real expiry logic lands in Feature 2.
0 * * * * expire_needs

# issue_delayed_token_rewards: grant once-only 24h milestone rewards for eligible resources and claims.
# Runs once per hour at minute 15.
15 * * * * issue_delayed_token_rewards

# process_resource_bid_notifications: warn about bids nearing resource expiry and automatically expire/cancel stale open bids.
# Runs once per hour at minute 17.
17 * * * * process_resource_bid_notifications

# issue_campaign_airdrop_coming_soon: notify eligible accounts when an approved campaign is within 48 hours of airdrop time.
# Runs once per hour at minute 18.
18 * * * * issue_campaign_airdrop_coming_soon

# issue_campaign_airdrop_payouts: pay approved campaign airdrops once per eligible account.
# Runs once per hour at minute 20.
20 * * * * issue_campaign_airdrop_payouts
