# backend/crontab
# Graphile Worker crontab
# Format: <cron expression> <task_name> [?options=json]
# Docs: https://worker.graphile.org/docs/cron
#
# expire_needs: mark campaign-needs past their target date as expired.
# Runs once per hour at minute 0. Real expiry logic lands in Feature 2.
0 * * * * expire_needs

# issue_delayed_token_rewards: grant once-only 24h milestone rewards for eligible resources and claims.
# Runs once per hour at minute 15.
15 * * * * issue_delayed_token_rewards

# process_resource_bid_notifications: warn about bids nearing resource expiry and automatically expire/cancel stale open bids.
# Runs once per hour at minute 17.
17 * * * * process_resource_bid_notifications

# issue_campaign_airdrop_coming_soon: notify eligible accounts when an approved campaign is within 48 hours of airdrop time.
# Runs once per hour at minute 18.
18 * * * * issue_campaign_airdrop_coming_soon

# issue_campaign_airdrop_payouts: pay approved campaign airdrops once per eligible account.
# Runs once per hour at minute 20.
20 * * * * issue_campaign_airdrop_payouts

# issue_notification_digests: queue the daily out-of-app digest email for accounts with due summary items.
# Runs once per day at 08:00 server time.
0 8 * * * issue_notification_digests

# deliver_mail_outbox: process queued emails from the unified mail outbox.
# Runs every minute to avoid noticeable delays in auth flows.
* * * * * deliver_mail_outbox

# deliver_push_notification_outbox: process queued realtime push notifications.
# Runs every minute to avoid noticeable delays in push delivery.
* * * * * deliver_push_notification_outbox
